# supervisord.conf
[supervisord]
nodaemon=true

[program:flask]
command=python3 /app/flask_app.py
autostart=true
autorestart=true
stdout_logfile=/app/logs/flask.log
stderr_logfile=/app/logs/flask_err.log

[program:telegram]
command=python3 /app/telegram_main.py
autostart=true
autorestart=true
stdout_logfile=/app/logs/telegram.log
stderr_logfile=/app/logs/telegram_err.log
